# contrib/package/community-profiles/files/etc/uci-defaults/freifunk-community-profiles
#!/bin/sh

rm -f /etc/config/profile_*
cd /usr/share/freifunk/community-profiles
for profile in profile_* ; do
  ln -s /usr/share/freifunk/community-profiles/${profile} /etc/config/${profile}
done
