# contrib/package/freifunk-common/files-common/etc/uci-defaults/freifunk
uci set uhttpd.main.rfc1918_filter=0
uci commit uhttpd

# has been converted to micrond
crontab -l | grep -v "/usr/sbin/ff_olsr_watchdog" | crontab -
