# contrib/package/freifunk-common/files-olsrv1/etc/uci-defaults/freifunk-olsrv1
#!/bin/sh

uci -m import freifunk <<EOF
config 'defaults' 'olsr_interfacedefaults'
	option 'Ip4Broadcast' '255.255.255.255'
EOF
